# systemd/system/dead-mans-switch.service
[Unit]
Description=Dead mans switch

[Service]
Type=oneshot
ExecStart=/usr/bin/env /var/www/html/dms.sh
RemainAfterExit=true
TimeoutStartSec=1800
