# systemd/system/dead-mans-switch.service
[Unit]
Description=Dead mans switch

[Service]
Type=oneshot
ExecStart=/bin/sh -c "/var/www/html/dms.sh"

[Install]
WantedBy=multi-user.target
