# systemd/system/dead-mans-switch.service
[Unit]
Description=Dead Man's Switch

[Service]
Type=forking
ExecStart=/bin/sh -c "/var/www/html/dms.sh"

[Install]
WantedBy=multi-user.target
