# systemd/dead-mans-switch.service
[Unit]
Description=Dead mans switch

[Service]
Type=simple
User=root
Group=root
ExecStart=/var/www/html/dms.sh
TimeoutStartSec=1800
